# DiskJockeyApplication/Resources/com.antimatterstudios.diskjockey.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.antimatterstudios.diskjockey.backend</string>
	<key>BundleProgram</key>
	<string>Contents/MacOS/diskjockey-backend</string>
	<key>ProgramArguments</key>
	<array>
		<string>diskjockey-backend</string>
		<string>--no-timeout</string>
	</array>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/diskjockey-backend-stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/diskjockey-backend-stderr.log</string>
</dict>
</plist>
